# src/etc/supervisor/conf.d/futel.conf
[program:hookswitch]
command=/opt/futel/src/hookswitch.py
;user=nobody
